# deployment/systemd-example/gemini-poc.service.example
# ~/.config/systemd/user/gemini-poc.service

[Unit]
Description=Gemini Chat PoC Pod Manager
# Start after network is ready and user login is complete
After=network-online.target default.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# IMPORTANT: Replace 'sheepdestroyer' with your actual username if different
#User=sheepdestroyer
WorkingDirectory=/mnt/LAB/gemini-chat-poc

# Execute the start script
ExecStart=/mnt/LAB/gemini-chat-poc/scripts/manage_pod.sh start

# Command to stop the pod
# Use the full path to podman found by 'which podman'
ExecStop=/usr/bin/podman pod stop gemini-poc-pod

# Optional: Set PATH if podman isn't found by default systemd path
# Environment="PATH=/bin:/usr/bin:/home/sheepdestroyer/bin:/home/sheepdestroyer/.local/bin"

[Install]
WantedBy=default.target
